# roles/lin-restic/templates/backup.service.j2
# ansible-managed
[Unit]
Description=Restic backup
[Service]
Type=oneshot
User=restic
Group=rclone
EnvironmentFile=/etc/restic/vars
ExecStart=/etc/restic/backup.sh
