# roles/lin-restic/templates/restic-backup.timer.j2
# ansible-managed
[Unit]
Description=Restic backup
[Timer]
OnCalendar={{ lin_restic_cron_time_backup }}
[Install]
WantedBy=timers.target
