# roles/lin-restic/templates/restic-prune.timer.j2
# ansible-managed
[Unit]
Description=Restic prune
[Timer]
OnCalendar={{ day }} {{ time }}
[Install]
WantedBy=timers.target
